# knowledge_base/product_config/celery-supervisor.conf
[program:celery]
command=/home/bo/projects/knowledge_base_project/.venv/bin/celery -A knowledge_base worker -l INFO -B
directory=/home/bo/projects/knowledge_base_project/knowledge_base/knowledge_base/user=bo
user=bo

numprocs=1

autostart=true
autorestart=true
startretries=3
startsecs=10
stopwaitsecs=10
stopasgroup=true
killasgroup=true
priority=998

stdout_logfile=/var/log/celery/base-nodes_worker.log
stderr_logfile=/var/log/celery/base-nodes_worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

environment=PYTHONPATH="/home/bo/projects/knowledge_base_project/knowledge_base/knowledge_base/",DJANGO_SETTINGS_MODULE="knowledge_base.settings",DJANGO_ENV="production"
# environment=DJANGO_SETTINGS_MODULE=knowledge_base.settings
